# etc/profile.d/hybris-ld.sh
#!/bin/sh

export HYBRIS_LD_LIBRARY_PATH=/usr/libexec/hal-droid/system/lib:/vendor/lib:/system/lib
